# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
GenericName[pl]=Dostosowanie temperatury barwowej
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
Name[pt_BR]=Redshift
Comment=Color temperature adjustment tool
Comment[ru]=Инструмент регулирования цветовой температуры
GenericName[cs]=Přizpůsobení barevné teploty
GenericName[hu]=Szín hőmérséklet beállítás
Comment[zh_CN]=色温调节工具
Terminal=false
Comment[it]=Strumento per la regolazione della temperatura del colore
Hidden=false
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
GenericName[lt]=Spalvos temperatūros reguliavimas
GenericName[de]=Farbtemperatur Anpassung
Comment[tr]=Renk sıcaklığı ayarlama aracı
Name[hu]=Redshift
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Name[it]=Redshift
Name[es]=Redshift
Name[fr]=Redshift
GenericName[it]=Regolazione della temperatura del colore
Comment[de]=Farbtemperatur Anpassungswerkzeug
Icon=redshift
Comment[es]=Herramienta para el ajuste de la temperatura de color
Name[zh_CN]=红移
Version=1.0
GenericName[tr]=Renk sıcaklığı ayarı
GenericName[da]=Justering af farvetemperatur
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
GenericName[zh_CN]=色温调节
Comment[hu]=Színhőmérséklet beállító eszköz
Exec=redshift-gtk
Name[ru]=Redshift
Name=Redshift
GenericName=Color temperature adjustment
X-GNOME-Autostart-enabled=true
Type=Application
GenericName[ru]=Настройка цветовой температуры
GenericName[fr]=Réglage de la température de couleur
Name[de]=Redshift
GenericName[nl]=Kleurtemperatuur aanpassing
GenericName[es]=Ajuste de la temperatura de color
StartupNotify=true
Comment[fr]=Outil de réglage de la température de couleur
Name[tr]=Redshift
Name[pl]=Redshift
Name[lt]=Redshift
Categories=Utility;
Comment[da]=Justeringsværktøj for farvetemperatur

